# rtux.service
[Unit]
Description=pressured — PSI-driven desktop responsiveness daemon
Documentation=https://github.com/justinstimatze/rtux
# NOTE: do NOT add `After=graphical.target` here. graphical.target is itself
# ordered After=multi-user.target, and this unit installs into
# multi-user.target (below) — the combination forms an ordering cycle that
# systemd breaks by DELETING our start job, so the daemon silently never runs
# at boot (it only came up when hand-started via install.sh). As a root cgroup
# daemon we don't need the graphical session to exist first; notifications
# already degrade gracefully until a user session is present.

[Service]
Type=simple
ExecStart=/usr/local/bin/pressured daemon
Restart=on-failure
RestartSec=5
# The protector must never be chosen as an OOM victim.
OOMScoreAdjust=-1000
# mlockall (CAP_IPC_LOCK) + raising our own limits if needed (CAP_SYS_RESOURCE).
AmbientCapabilities=CAP_IPC_LOCK CAP_SYS_RESOURCE
# Keep our own footprint bounded — the daemon stays tiny.
MemoryMax=64M

# --- Sandboxing ---
# This is a root daemon, so it's hardened as far as its job allows. It is NOT
# heavily confinable: it must write cgroup v2 controls all over /sys/fs/cgroup,
# read every process's /proc, write other processes' oom_score_adj, and drop to
# the graphical user (runuser) to send notifications. So the directives that would
# break those — ProtectControlGroups=yes, ProtectKernelTunables=yes,
# ProtectProc=/ProcSubset=, ProtectHome=, PrivateTmp= — are deliberately left off.
# Everything below is safe with the above and still meaningfully shrinks the
# blast radius.
#
# CAP_SYS_PTRACE is required to READ (never to trace): the ptrace-gated parts of
# /proc — notably /proc/<pid>/cwd — are what let us name a paused/killed session
# by its working directory ("claude · rtux" rather than a bare "claude"). Being
# uid 0 is NOT sufficient: ptrace_may_access() only skips its capability check
# when the reader's creds MATCH the target's, so root reading a uid-1000 process
# still needs CAP_SYS_PTRACE and otherwise gets EACCES. Without it every session
# silently degraded to an unnamed "claude" — which is precisely the legibility
# (which session died? can I resume it?) this daemon exists to provide. The
# world-readable /proc/<pid>/comm kept working throughout, which is why the
# breakage looked like a naming regression rather than a permissions one.
CapabilityBoundingSet=CAP_IPC_LOCK CAP_SYS_RESOURCE CAP_SETUID CAP_SETGID CAP_CHOWN CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH CAP_SYS_PTRACE
NoNewPrivileges=true
ProtectSystem=full
ProtectControlGroups=false
ProtectClock=true
ProtectHostname=true
ProtectKernelModules=true
ProtectKernelLogs=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_NETLINK
# Filtered syscalls return EPERM rather than killing the daemon — a filter miss
# degrades gracefully (e.g. a logged mlockall warning) instead of taking the
# protector down.
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
